# files/etc/hotplug.d/iface/50-gluon-wan-dnsmasq
if [ "$INTERFACE" = 'umts' -o "$INTERFACE" = 'tether' -o "$INTERFACE" = 'wan' -o "$INTERFACE" = 'wan6' ]; then
   /lib/gluon/wan-dnsmasq/update.lua
fi
